# docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=info

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true

[program:php-fpm83]
command=php-fpm83
autostart=true
autorestart=true
